# packages/fujin-cli/fujin_cli-0.2.0.tar.gz/fujin_cli-0.2.0/src/fujin/templates/web.socket
[Unit]
Description={app} socket

[Socket]
ListenStream=/run/{app}.sock
SocketUser=www-data
SocketGroup=www-data
SocketMode=0660

[Install]
WantedBy=sockets.target
