# packages/fujin-cli/fujin_cli-0.2.0-py3-none-any.whl/fujin/templates/other.service
[Unit]
Description={app} Worker

[Service]
User={user}
Group=www-data
WorkingDirectory={project_dir}
ExecStart={command}
EnvironmentFile={project_dir}/.env
Restart=always
#StandardOutput=append:/var/log/your_project/qcluster.out.log
#StandardError=append:/var/log/your_project/qcluster.err.log

[Install]
WantedBy=multi-user.target